# tests/containerfiles/lbi/usr/share/containers/systemd/curl.container
[Container]
Image=quay.io/curl/curl:latest
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
