# systemd/bootupd.service
[Unit]
Description=bootloader update daemon

[Service]
Type=notify
ExecStart=/usr/libexec/bootupd daemon
