# systemd/bootc-status-updated.target
[Unit]
Description=Target for bootc status changes
Documentation=man:bootc-status-updated.target(8)
StopWhenUnneeded=true
