# systemd/bootupd.socket
[Socket]
ListenSequentialPacket=/run/bootupd.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
