# systemd/bootc-status-updated.path
[Unit]
Description=Monitor bootc for status changes
Documentation=man:bootc-status-updated.path(8)

[Path]
PathChanged=/ostree/bootc
Unit=bootc-status-updated.target

[Install]
WantedBy=multi-user.target
