# systemd/bootupd.service
[Unit]
Description=bootloader update daemon
Documentation=https://github.com/coreos/bootupd

[Service]
Type=notify
Environment=BOOTUPD_VERBOSITY="-v"
ExecStart=/usr/libexec/bootupd daemon $BOOTUPD_VERBOSITY
# On general principle
ProtectHome=yes
# So we can remount /boot writable
MountFlags=slave
